# charybdis-atheme/supervisord.conf
[supervisord]
nodaemon=true

[program:ircd]
command=/ircd/bin/ircd -foreground
user=ircd

[program:atheme]
command=/atheme/bin/atheme-services -n
user=ircd
